# raid-check-serial.service
[Unit]
Description=RAID check with media-aware concurrency limits

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/raid-check-serial
ExecStart=/usr/local/sbin/raid-check-serial.sh
